# senpan.service
[Unit]
Description=Senpan App Suite Service
Documentation=https://github.com/jwill89
After=network.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/senpan
# -webroot must point at the Apache DocumentRoot so uploaded raffle images are
# written to <DocumentRoot>/images/raffles/ and served by Apache. Adjust the
# path to match your server (see deploy/README.md).
ExecStart=/opt/senpan/app-suite -addr :8080 -db /opt/senpan/data/database.sqlite -webroot /var/www/apps.senpan.cafe

# Auto-restart on failure, with increasing back-off
Restart=on-failure
RestartSec=5

# Environment
#Environment=APPSUITE_SESSION_SECRET=changeme
# Only if cross-origin browser clients need the API (normally unset — same-origin):
#Environment=APPSUITE_CORS_ORIGINS=https://example.com

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
# The app needs write access to its own dir (DB) and to the document root's
# images/ folder where raffle uploads are written. Adjust the second path to
# match your -webroot (only the images/ subtree needs to be writable).
ReadWritePaths=/opt/senpan /var/www/apps.senpan.cafe/images

# Logging (stdout/stderr → journalctl)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=senpan

[Install]
WantedBy=multi-user.target

