# deploy/senpan.service
# systemd unit TEMPLATE for the Senpan App Suite backend. This tracked copy keeps
# the concrete host layout OUT of the repo: copy it to
# /etc/systemd/system/<service-name>.service on the host and replace the
# placeholders below — <service-user>, /opt/<app>, <webroot>, <service-name> —
# with your values. See deploy/README.md for a worked example. (deploy.ps1 does
# NOT ship this file; the running unit lives only on the host.)

[Unit]
Description=Senpan App Suite Service
Documentation=https://github.com/jwill89
After=network.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
# Run as a dedicated, unprivileged system account — nothing here needs root
# (the listen addr is :8080, non-privileged; all paths are ordinary files). See
# deploy/README.md "Running as a non-root user" for the one-time useradd + chown.
User=<service-user>
Group=<service-user>
WorkingDirectory=/opt/<app>
# -webroot must point at the Apache DocumentRoot so uploaded raffle images are
# written to <webroot>/images/raffles/ and served by Apache.
ExecStart=/opt/<app>/app-suite -addr :8080 -db /opt/<app>/data/database.sqlite -webroot <webroot>

# Auto-restart on failure, with increasing back-off
Restart=on-failure
RestartSec=5

# ── Configuration & secrets ──
# Config comes from the flags in ExecStart above OR from APPSUITE_* env vars. Two
# of them are SECRETS — the session cookie secret and the Turnstile secret — so
# keep them OUT of this unit file: an inline `Environment=SECRET=…` value is
# world-readable via `systemctl show` and /proc/<pid>/environ. Put secrets in a
# root-owned 0600 env file instead and load it here (the leading `-` makes it
# optional, so the service still starts if it's absent):
#
#   /etc/senpan/senpan.env   (chown root:root, chmod 600):
#     APPSUITE_SESSION_SECRET=<fixed 64+ hex, e.g. `openssl rand -hex 32`; a fixed
#                              value keeps sessions alive across restarts>
#     APPSUITE_TURNSTILE_SECRET=<Cloudflare Turnstile secret key>
#     APPSUITE_TURNSTILE_SITEKEY=<Cloudflare Turnstile public site key>  # public, but
#                              kept alongside the pair for convenience
#
#EnvironmentFile=-/etc/senpan/senpan.env
#
# Non-secret values are fine inline. Only if cross-origin browser clients need the
# API (normally unset — same-origin):
#Environment=APPSUITE_CORS_ORIGINS=https://example.com

# Security hardening. These sandbox the service so a compromise of the process
# has minimal reach: no privilege escalation, a read-only filesystem except the
# explicit ReadWritePaths below, no access to /home, a private /tmp and /dev, and
# a locked-down kernel/namespace/syscall surface. Adjust only with matching app
# changes — every writable path the app touches must be in ReadWritePaths.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
# The service only speaks TCP (listen :8080) and may use local sockets for DNS.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
# A Go network daemon needs nothing outside the standard service syscall set.
SystemCallFilter=@system-service
# Reject any non-native syscall ABI so the filters above (and W^X below) can't be
# bypassed through a secondary architecture's calling convention.
SystemCallArchitectures=native
# W^X: forbid writable-executable mappings (blocks classic code-injection). Safe
# here because nothing in this binary JITs — it's pure Go with no cgo, and the
# SQLite driver (ncruces/go-sqlite3 v0.35+) transpiles SQLite's WebAssembly to
# native Go ahead of time via wasm2go (see go-sqlite3-wasm: "Code generated by
# wasm2go"), so SQLite runs from the read-only .text segment with no runtime code
# generation. (This was NOT safe under older ncruces, which JITed WASM via
# wazero's compiler — if the driver is ever pinned back to a wazero-JIT version,
# this line must be removed again.)
MemoryDenyWriteExecute=yes

# The app needs write access to its own dir (DB + SQLite -wal/-shm) and to the
# document-root subtrees it writes into: images/ (raffle/announcement/category
# uploads), fonts/ (uploaded fonts + the fonts/.woff2 conversion cache), and
# carrd/ (per-project image folders). Every subtree the app writes MUST appear
# here or writes fail with EACCES under ProtectSystem=strict. These paths MUST
# match the deployed -webroot in ExecStart above (same <webroot> value).
ReadWritePaths=/opt/<app> <webroot>/images <webroot>/fonts <webroot>/carrd
# Creates + owns /var/log/senpan (service user, 0755) and exposes it as
# $LOGS_DIRECTORY. Provides a writable location for a file log under
# ProtectSystem=strict; journald still captures stdout/stderr below.
LogsDirectory=senpan

# Logging (stdout/stderr → journalctl)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=<service-name>

[Install]
WantedBy=multi-user.target
